# rc.local
#!/bin/sh

if [ -f /etc/ipset.conf ]; then
    ipset restore -file /etc/ipset.conf
fi

if [ -f /etc/iptables.rules ]; then
    iptables-restore < /etc/iptables.rules
fi

ss-redir -c /etc/shadowsocks-libev/config.json -f /var/run/shadowsocks.pid

exit 0